# nonebot.conf
[program:nonebot]
directory=/koyeb-nb2
command=sh -c ./run-uvicorn.sh
autostart=true
autorestart=false
stderr_logfile=/koyeb-nb2/noenbot.err.log
stdout_logfile=/koyeb-nb2/nonebot.out.log
